# infra/systemd/agent-monitor.service
[Unit]
Description=Agent Monitor for GCP Google Chat Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=oracle
Group=oinstall
WorkingDirectory=/u01/app/agent-monitor/agent
EnvironmentFile=/u01/app/agent-monitor/.env.runtime
ExecStart=/u01/app/agent-monitor/venv/bin/uvicorn main:app --host ${AGENT_HOST} --port ${AGENT_PORT} --log-level info --no-access-log
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
